# products/qestro/scripts/production/config/production/qestro-health-monitor.service
[Unit]
Description=Qestro Platform Health Monitor
After=network.target
Wants=network-online

[Service]
Type=oneshot
ExecStart=/opt/qestro/scripts/production/health-check.sh
StandardOutput=journal:syslog
StandardError=journal:syslog
User=root
Group=root

[Install]
WantedBy=multi-user.target
